# scripts/crontab
0 2 * * * . /etc/environment && /bin/bash /scripts/backup.sh >> /proc/1/fd/1 2>&1
